# root/etc/services.d/radarr/run
#!/usr/bin/with-contenv bash

redirect_cmd() {
	if [[ "${DEBUG}" == "true" ]]; then
		"$@"
	else
		"$@" >/dev/null 2>&1
	fi
}

echo "Starting Radarr."
redirect_cmd exec s6-setuidgid abc \
	/app/radarr/bin/Radarr \
	-nobrowser \
	-data=/config
